# bananabot.service
[Unit]
Description=BananaBot Discord Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/bananabot
Environment="PATH=/root/bananabot/venv/bin"
ExecStart=/root/bananabot/venv/bin/python /root/bananabot/start.py
Restart=always
RestartSec=10
StandardOutput=append:/root/bananabot/bot.log
StandardError=append:/root/bananabot/bot.log

[Install]
WantedBy=multi-user.target